# scripts/com.parthchandak.litreview-healthwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.parthchandak.litreview-healthwatch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@@WATCHDOG_SCRIPT@@</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>@@STATE_DIR@@/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>@@STATE_DIR@@/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@@LAUNCHD_PATH@@</string>
        <key>PM2_BIN</key>
        <string>@@PM2_BIN@@</string>
    </dict>
</dict>
</plist>
